# system/99-auto-power-profile.rules
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", ACTION=="change", RUN+="/usr/bin/systemctl --no-block start auto-power-profile.service"
